# modules/rails/templates/upstart-sidekiq.conf.erb
start on starting <%= @app_name %>
stop on stopping <%= @app_name %>
respawn

exec su - <%= @app_name %> -c 'cd /u/apps/<%= @app_name %>/current && RAILS_ENV=<%= @rails_env %> exec ruby1.9.1 -S bundle exec sidekiq -P /u/apps/<%= @app_name %>/shared/pids/sidekiq.pid >> /u/apps/<%= @app_name %>/shared/log/sidekiq.log 2>&1'
