# modules/rails/templates/upstart-thin-multi.conf.erb
start on starting <%= @app_name %>-thin
stop on stopping <%= @app_name %>-thin
respawn

exec su - <%= @app_name %> -c 'cd /u/apps/<%= @app_name %>/current && RAILS_ENV=<%= @rails_env %> exec <%= @bundle_exec %> thin -e <%= @rails_env %> -S /u/apps/<%= @app_name %>/shared/thin.<%= @app_server_id %>.sock start >> /u/apps/<%= @app_name %>/shared/log/thin.<%= @app_server_id %>.log 2>&1'
