# modules/rails/templates/upstart-sidekiq.conf.erb
start on starting <%= @app_name %>
stop on stopping <%= @app_name %>
respawn

env RAILS_ENV="<%= @rails_env %>"
export RAILS_ENV
setuid <%= @app_name %>
setgid <%= @app_name %>
chdir /u/apps/<%= @app_name %>/current

exec <%= @bundle_exec %> sidekiq -P /u/apps/<%= @app_name %>/shared/pids/sidekiq.pid
